# minecraft.service
[Unit]
Description=Minecraft Server
After=network.target

[Service]
Type=simple
User=deploy
Group=deploy
WorkingDirectory=/home/deploy/minecraft-server

# RAM configuration
Environment="MIN_RAM=2G"
Environment="MAX_RAM=4G"

# Start and stop commands
ExecStart=/usr/bin/java -Xms${MIN_RAM} -Xmx${MAX_RAM} -jar server.jar nogui
ExecStop=/bin/kill -SIGINT $MAINPID
Restart=on-failure
RestartSec=10

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/home/deploy/minecraft-server

[Install]
WantedBy=multi-user.target
