# minecraft.service
[Unit]
Description=Minecraft Server
After=network.target

[Service]
Type=simple
User=deploy
Group=deploy
WorkingDirectory=/home/deploy/minecraft-server
ExecStart=/usr/bin/java -Xms2G -Xmx4G -jar server.jar nogui
ExecStop=/bin/kill -SIGINT $MAINPID
Restart=on-failure
RestartSec=10

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/home/deploy/minecraft-server

[Install]
WantedBy=multi-user.target
