# minecraft.service
[Unit]
Description=Minecraft Server
After=network.target

[Service]
Type=simple
User=deploy
Group=deploy
WorkingDirectory=/home/deploy/minecraft-server
# Load configuration and start server
ExecStart=/bin/bash -c 'source /home/deploy/minecraft-server/config.sh 2>/dev/null || { MIN_RAM=2G; MAX_RAM=4G; }; exec /usr/bin/java -Xms${MIN_RAM} -Xmx${MAX_RAM} -jar server.jar nogui'
ExecStop=/bin/kill -SIGINT $MAINPID
Restart=on-failure
RestartSec=10

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/home/deploy/minecraft-server

[Install]
WantedBy=multi-user.target
